# install/conf/snmpd.ini
[program:snmpd]
user=root
directory=/
command=/usr/sbin/snmpd -c /etc/snmp/snmpd.conf
autostart=true
autorestart=true
startsecs=5
stopsignal=HUP
stopasgroup=true
stopwaitsecs=5
stdout_logfile_maxbytes=20MB
stdout_logfile=/var/log/supervisor/%(program_name)s-out.log
stderr_logfile_maxbytes=20MB
stderr_logfile=/var/log/supervisor/%(program_name)s-err.log